# base/resources/wallpaper/wallpaper-systemd.path
[Unit]
Description=Refresh wallpaper when file changes

[Path]
PathModified=/opt/pi4j-os/wallpaper-dynamic.jpg
