# image/resources/wallpaper-systemd.path
[Unit]
Description=Refresh wallpaper when file changes

[Path]
PathModified=/opt/fhnw/wallpaper-dynamic.jpg
